# deployment/sandbox-controller/nubian-screen-recorder.service
[Unit]
Description=Nubian Rolling Screen Recorder
After=graphical-session.target nubian-universal-computer-agent.service

[Service]
Type=simple
WorkingDirectory=%h/.local/share/nubian-universal-computer-agent
Environment=DISPLAY=:0
Environment=XAUTHORITY=/run/user/1000/gdm/Xauthority
Environment=XDG_SESSION_TYPE=x11
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
Environment=SCREEN_RECORDER_DIR=/logs/screen-recorder
Environment=SCREEN_RECORDER_FPS=5
Environment=SCREEN_RECORDER_RETENTION_SEC=300
Environment=SCREEN_RECORDER_PNG_COMPRESS_LEVEL=1
ExecStart=/usr/bin/python3 %h/.local/share/nubian-universal-computer-agent/screen_recorder.py
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
